# ansible/templates/quicknotes-ansible-pull.timer.j2
[Unit]
Description=Run QuickNotes ansible-pull every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
Unit={{ quicknotes_pull_service_name }}

[Install]
WantedBy=timers.target